# newt/rootfs/etc/s6-overlay/s6-rc.d/newt/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2155

# Declare variables
export PANGOLIN_ENDPOINT=$(bashio::config 'pangolin_endpoint')
export NEWT_ID=$(bashio::config 'newt_id')
export NEWT_SECRET=$(bashio::config 'newt_secret')
export MTU=$(bashio::config 'mtu')
export DNS=$(bashio::config 'dns')
export ACCEPT_CLIENTS=$(bashio::config 'ACCEPT_CLIENTS')
export PING_INTERVAL=$(bashio::config 'ping_interval')
export PING_TIMEOUT=$(bashio::config 'ping_timeout')
export LOG_LEVEL=$(bashio::config 'log_level')

# Launch Newt
exec /usr/local/bin/newt
