# newt/rootfs/etc/s6-overlay/s6-rc.d/newt/run
#!/usr/bin/with-contenv bashio

# Declare variables
PANGOLIN_ENDPOINT=$(bashio::config 'pangolin_endpoint')
NEWT_ID=$(bashio::config 'newt_id')
NEWT_SECRET=$(bashio::config 'newt_secret')

# Launch Newt
export PANGOLIN_ENDPOINT NEWT_ID NEWT_SECRET
exec /usr/local/bin/newt
